# task1/monitoring.service
[Unit]
Description=Мониторинг веб-приложения
After=network.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/monitoring.sh
